# systemd-unit/credfeto-package-cache.service
[Unit]
Description=Update Credfeto Package Cache
After=network.target

[Service]
User=markr
Type=oneshot
IOSchedulingClass=idle
CPUSchedulingPolicy=idle
WorkingDirectory=/home/markr/work/credfeto-dev-package-cache
ExecStart=/bin/flock -n /tmp/devpkg.lock -c /home/markr/work/credfeto-dev-package-cache/update

